# metro-bistro-backup.timer
[Unit]
Description=Run Metro Bistro backup daily at 2 AM
Requires=metro-bistro-backup.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target 